# customization/crontab/crontab.sh
0 7 * * 1-5 ~/scripts/automate_initialization/automate_initialization.sh
0 12 * * 1-5 ~/scripts/automate_finances/query_balance.sh
